# srcs/phpmyadmin/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=sh -c "nginx -g 'daemon off;' && kill $(cat supervisord.pid)"

[program:php-fpm]
command=sh -c "/usr/sbin/php-fpm7 --nodaemonize && kill $(cat supervisord.pid)"